# templates/systemd/letsencrypt.renew.timer.j2
[Unit]
Description=Run daily letsencrypt renew

[Timer]
OnCalendar = daily
Persistent=true

[Install]
WantedBy=timers.target
