# templates/renew.cron.j2
# m  h          d M W
  30 5          * * *           root    /bin/bash --login {{renew_script_path}} >> {{renew_logfile}}
